# init/timer.service
[Unit]
Description=Run Daily Backup

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 /home/arslan/projects/operating-systems-course/init/users.py

